# data/learning.json
[
  {"id": "process", "title": "Windows Process", "category": "OS Core",
   "summary": "A running instance of a program with its own virtual address space, handles and threads.",
   "purpose": "Processes are the unit of execution and isolation in Windows. Each has a PID, a parent (PPID), a token defining its privileges, and an integrity level.",
   "architecture": "A process owns a private virtual address space, one or more threads, loaded modules (DLLs), and handles to kernel objects. Creation flows through CreateProcess -> kernel process/thread objects.",
   "dependencies": "ntdll.dll and the loader map dependencies; the token comes from the parent or an explicit logon.",
   "windows_logs": "Security 4688 (process creation, with command line if auditing is enabled) and 4689 (process exit).",
   "linux_logs": "auditd execve records; /var/log/auth.log for privilege changes.",
   "sysmon": ["Event ID 1: Process Create (hashes, full command line, parent image)", "Event ID 5: Process Terminated"],
   "common_attacks": ["Masquerading as a system process (T1036)", "Process injection (T1055)", "Parent PID spoofing (T1134.004)"],
   "detection": "Baseline parent-child relationships, alert when Office apps spawn shells, and verify image paths and signatures against expected locations.",
   "sigma": "title: Office Spawning PowerShell\nlogsource:\n  product: windows\n  category: process_creation\ndetection:\n  sel:\n    ParentImage|endswith: '\\WINWORD.EXE'\n    Image|endswith: '\\powershell.exe'\n  condition: sel",
   "yara": "",
   "hardening": "Enable command-line auditing and Sysmon; apply Attack Surface Reduction rules to block Office child processes.",
   "enterprise": "Process telemetry is the backbone of EDR detections and endpoint threat hunting."},

  {"id": "service", "title": "Windows Service", "category": "OS Core",
   "summary": "A background program managed by the Service Control Manager, often running as SYSTEM.",
   "purpose": "Services provide always-on functionality (networking, printing, security) and start automatically without a user logged on.",
   "architecture": "The SCM (services.exe) reads service config from the registry (Start, ImagePath, ObjectName) and hosts many services inside shared svchost.exe groups.",
   "dependencies": "Most services depend on RpcSs/DcomLaunch; config lives under HKLM\\SYSTEM\\CurrentControlSet\\Services.",
   "windows_logs": "System 7045 (new service installed), 7036 (state change); Security 4697 (service installed).",
   "linux_logs": "systemd journal (journalctl); unit files under /etc/systemd/system.",
   "sysmon": ["Event ID 1 for the hosting process", "Event ID 13 for service registry changes"],
   "common_attacks": ["New malicious service for persistence (T1543.003)", "Service binary hijack", "Unquoted service path abuse"],
   "detection": "Alert on 7045 with binaries outside System32, services running from user-writable paths, and svchost with an unexpected parent.",
   "sigma": "title: Suspicious New Service Binary\nlogsource:\n  product: windows\n  service: system\ndetection:\n  sel:\n    EventID: 7045\n    ImagePath|contains: '\\Users\\'\n  condition: sel",
   "yara": "",
   "hardening": "Restrict service creation rights, fix unquoted paths, and monitor the Services registry hive.",
   "enterprise": "Autoruns and service baselining are core to persistence hunting."},

  {"id": "driver", "title": "Kernel Driver", "category": "Kernel",
   "summary": "Kernel-mode code (.sys) running at Ring 0 to talk to hardware or extend the OS.",
   "purpose": "Drivers bridge hardware and the OS and can filter file, registry and network operations (used by both EDR and malware).",
   "architecture": "Drivers load per their Start type (Boot/System/Auto/Demand) and register with the I/O manager; Driver Signature Enforcement validates signatures.",
   "dependencies": "Most modern drivers build on the Kernel-Mode Driver Framework (Wdf01000.sys).",
   "windows_logs": "System 7026 (boot driver load failure); Code Integrity 3001/3033 (unsigned/blocked driver).",
   "linux_logs": "dmesg / kernel ring buffer; modprobe and lsmod for loaded modules.",
   "sysmon": ["Event ID 6: Driver Loaded (image, signature status)"],
   "common_attacks": ["Bring Your Own Vulnerable Driver (T1068)", "Kernel rootkits hooking SSDT/callbacks", "Disabling DSE via test-signing"],
   "detection": "Compare loaded drivers against the Microsoft vulnerable-driver blocklist; alert on unsigned or newly loaded kernel modules.",
   "sigma": "title: Vulnerable Driver Loaded\nlogsource:\n  product: windows\n  service: sysmon\ndetection:\n  sel:\n    EventID: 6\n    Signature|contains: 'unknown'\n  condition: sel",
   "yara": "",
   "hardening": "Enable HVCI/Memory Integrity, the Microsoft vulnerable-driver blocklist, and kernel lockdown.",
   "enterprise": "BYOVD is a leading path to kernel compromise; driver telemetry is essential."},

  {"id": "registry", "title": "Windows Registry", "category": "Configuration",
   "summary": "A hierarchical database of keys and values storing OS and application configuration.",
   "purpose": "The registry holds boot, service, security and per-user settings - and, as a result, most persistence real estate.",
   "architecture": "Organised into hives (HKLM, HKCU, ...) backed by files on disk. Keys contain typed values (REG_SZ, REG_DWORD, etc.).",
   "dependencies": "The SYSTEM hive is read during boot before most drivers and services.",
   "windows_logs": "Security 4657 (registry value modified) when object auditing is enabled.",
   "linux_logs": "N/A (Linux uses config files under /etc).",
   "sysmon": ["Event ID 12: Registry object create/delete", "Event ID 13: Registry value set"],
   "common_attacks": ["Run key persistence (T1547.001)", "Winlogon Shell/Userinit hijack (T1547.004)", "Image File Execution Options debugger (T1546.012)"],
   "detection": "Monitor autostart locations and Winlogon keys with Sysmon 12/13; alert on Run values pointing to AppData or temp.",
   "sigma": "title: Run Key Persistence\nlogsource:\n  product: windows\n  service: sysmon\ndetection:\n  sel:\n    EventID: 13\n    TargetObject|contains: '\\CurrentVersion\\Run'\n  condition: sel",
   "yara": "",
   "hardening": "Restrict write access to autostart keys and forward Sysmon registry events to the SIEM.",
   "enterprise": "Autoruns triage and registry baselining are standard IR steps."},

  {"id": "lsass", "title": "LSASS & Credentials", "category": "Security",
   "summary": "The Local Security Authority Subsystem enforces security policy and caches credential material.",
   "purpose": "lsass.exe validates logons and holds secrets (NTLM hashes, Kerberos tickets) in memory - the prime credential-theft target.",
   "architecture": "Runs as SYSTEM from wininit.exe; authentication packages (msv1_0, kerberos) plug into it.",
   "dependencies": "Backed by the SAM/SECURITY hives and, in a domain, the domain controller.",
   "windows_logs": "Security 4624/4625 (logon success/failure), 4768/4769 (Kerberos), 4776 (NTLM).",
   "linux_logs": "N/A (Linux uses PAM + /etc/shadow).",
   "sysmon": ["Event ID 10: ProcessAccess to lsass.exe (watch GrantedAccess 0x1010/0x1410)"],
   "common_attacks": ["Credential dumping via MiniDump/comsvcs.dll (T1003.001)", "Mimikatz sekurlsa", "Pass-the-Hash / Pass-the-Ticket"],
   "detection": "Alert on any non-system process opening lsass with VM_READ; watch for lsass.dmp file creation.",
   "sigma": "title: LSASS Memory Access\nlogsource:\n  product: windows\n  service: sysmon\ndetection:\n  sel:\n    EventID: 10\n    TargetImage|endswith: '\\lsass.exe'\n    GrantedAccess: '0x1410'\n  condition: sel",
   "yara": "rule Mimikatz_Strings { strings: $a = \"sekurlsa\" ascii wide condition: $a }",
   "hardening": "Enable Credential Guard, RunAsPPL (LSASS protection), and restrict debug privilege.",
   "enterprise": "LSASS protection and detection are among the highest-value defensive controls."},

  {"id": "powershell", "title": "PowerShell", "category": "Execution",
   "summary": "A powerful scripting and automation shell, and a favourite living-off-the-land tool for attackers.",
   "purpose": "PowerShell administers Windows at scale; its .NET access and in-memory execution make it ideal for fileless attacks.",
   "architecture": "The engine hosts scriptblocks; logging is provided by module, scriptblock and transcription logging.",
   "dependencies": "Relies on .NET; can call Win32 APIs directly.",
   "windows_logs": "PowerShell/Operational 4103 (module) and 4104 (scriptblock); Security 4688 for process start.",
   "linux_logs": "PowerShell Core logs to syslog when configured.",
   "sysmon": ["Event ID 1 with the decoded command line", "Event ID 3 for outbound download cradles"],
   "common_attacks": ["Encoded commands (-enc) (T1059.001)", "Download cradles (IEX DownloadString)", "In-memory .NET / reflective loading"],
   "detection": "Enable scriptblock logging (4104) and alert on -enc, hidden windows, DownloadString and AMSI bypass strings.",
   "sigma": "title: Encoded PowerShell\nlogsource:\n  product: windows\n  category: process_creation\ndetection:\n  sel:\n    Image|endswith: '\\powershell.exe'\n    CommandLine|contains: ' -enc '\n  condition: sel",
   "yara": "",
   "hardening": "Enable Constrained Language Mode, scriptblock logging, AMSI, and application control (WDAC).",
   "enterprise": "Scriptblock logging is one of the highest-signal Windows telemetry sources."},

  {"id": "kerberos", "title": "Kerberos", "category": "Authentication",
   "summary": "The default domain authentication protocol, based on tickets issued by a Key Distribution Center.",
   "purpose": "Kerberos provides single sign-on using time-limited, encrypted tickets instead of transmitting passwords.",
   "architecture": "The client gets a TGT (AS exchange), then service tickets (TGS exchange), then presents them to services (AP exchange).",
   "dependencies": "Requires synchronised clocks (within 5 minutes) and DNS; the krbtgt account key encrypts all TGTs.",
   "windows_logs": "Security 4768 (TGT requested), 4769 (service ticket), 4771 (pre-auth failed).",
   "linux_logs": "krb5 logs; SSSD logs on domain-joined Linux.",
   "sysmon": ["Correlate 4769 bursts with the requesting process"],
   "common_attacks": ["Kerberoasting (T1558.003)", "Golden Ticket (T1558.001)", "Silver Ticket / Pass-the-Ticket"],
   "detection": "Alert on many 4769 RC4 (0x17) tickets from one account and on TGTs with anomalous lifetimes.",
   "sigma": "title: Potential Kerberoasting\nlogsource:\n  product: windows\n  service: security\ndetection:\n  sel:\n    EventID: 4769\n    TicketEncryptionType: '0x17'\n  condition: sel | count() by AccountName > 10",
   "yara": "",
   "hardening": "Use strong, long service-account passwords (or gMSA), disable RC4, and protect the krbtgt account.",
   "enterprise": "Kerberos abuse is central to Active Directory attack paths."},

  {"id": "sysmon", "title": "Sysmon & Event Logs", "category": "Detection",
   "summary": "System Monitor (Sysmon) enriches Windows logging with high-fidelity endpoint telemetry.",
   "purpose": "Sysmon records process, network, image-load and registry activity with hashes and full command lines - the raw material for detection.",
   "architecture": "A driver + service writing to the Microsoft-Windows-Sysmon/Operational log, driven by an XML config that filters noise.",
   "dependencies": "Best paired with a log forwarder (WEF) and a SIEM.",
   "windows_logs": "Security 4688/4624/4672; System 7045; the dedicated Sysmon/Operational channel.",
   "linux_logs": "Sysmon for Linux writes to syslog; auditd is the native equivalent.",
   "sysmon": ["1 Process Create", "3 Network", "7 Image Load", "8 CreateRemoteThread", "10 ProcessAccess", "11 FileCreate", "13 Registry Set", "22 DNS Query"],
   "common_attacks": ["Log clearing (T1070.001, Security 1102)", "Disabling Sysmon/EventLog service"],
   "detection": "Alert on Security 1102 (log cleared) and on stops of the EventLog or Sysmon service.",
   "sigma": "title: Windows Event Log Cleared\nlogsource:\n  product: windows\n  service: security\ndetection:\n  sel:\n    EventID: 1102\n  condition: sel",
   "yara": "",
   "hardening": "Deploy a tuned Sysmon config, forward logs off-host, and alert on logging tampering.",
   "enterprise": "A good Sysmon baseline is the single biggest boost to endpoint detection coverage."}
]
